# packages/rsidatasciencetools/rsidatasciencetools-0.0.5-py3-none-any.whl/rsidatasciencetools/dockerutils/supervisord.conf
[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/log/supervisord.pid                ; pidfile location
nodaemon=True                                   ; do not run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port=*:9001
username=mlapiuser
password=ml@piu53r

[supervisorctl]
user = root
username = mlapiuser
password = ml@piu53r
prompt = sprvsrd

# [program:pythontrainmodel]
# stdout_logfile_maxbytes=50MB
# numprocs=1
# user=apiuser
# autostart=false
# autorestart=false
# startsecs=15
# startretries=0
# stopwaitsecs=15
# exitcodes=0,2
# killasgroup=true
# stopsignal=TERM
# stopasgroup=true
# directory=/App/formsML/mlmodels
# command=bash -c 'source /App/docker_build.env && python train.py --data-dir=$RSI_TRAINDATADIR -d; python -c "import time; time.sleep(10)"'
# environment=RSI_TRAINDATADIR='%(ENV_RSI_TRAINDATADIR)s',
#     RSI_ML_API_PORT='%(ENV_RSI_ML_API_PORT)s'

[program:pythonapi]
stdout_logfile_maxbytes=50MB
numprocs=1
user=apiuser
autostart=true
autorestart=unexpected
startsecs=15
startretries=3
stopwaitsecs=15
exitcodes=0,2
killasgroup=true
stopsignal=TERM
stopasgroup=true
directory=/App/formsML/mlapi
command=bash -c 'source /App/docker_build.env && uvicorn --host 0.0.0.0 --port $RSI_ML_API_PORT --reload forms:api'
environment=RSI_TRAINDATADIR='%(ENV_RSI_TRAINDATADIR)s',
    RSI_ML_API_PORT='%(ENV_RSI_ML_API_PORT)s'

[program:jupyternb]
stdout_logfile_maxbytes=50MB
numprocs=1
user=apiuser
autostart=false
autorestart=unexpected
startsecs=15
startretries=2
stopwaitsecs=15
exitcodes=0,2
killasgroup=true
stopsignal=TERM
stopasgroup=true
directory=/home/apiuser/notebooks
command=jupyter notebook -y --ip=0.0.0.0 --no-browser 
# --notebook-dir=/home/apiuser/notebooks/
environment=RSI_TRAINDATADIR='%(ENV_RSI_TRAINDATADIR)s',
    RSI_ML_API_PORT='%(ENV_RSI_ML_API_PORT)s',
    RSI_JUPYTER_PORT='%(ENV_RSI_JUPYTER_PORT)s',
    HOME=/home/apiuser/,
    USER=apiuser
