# packages/django-channel-tasks/django-channel-tasks-0.5.0.dev6.tar.gz/django-channel-tasks-0.5.0.dev6/services/channel-tasks.template.service
[Unit]
Description=channel-tasks

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=${CHANNEL_TASKS_HOME}/bin/setup-channel-tasks.sh
Requires=redis.service
Requires=unit.service
Environment="CHANNEL_TASKS_HOME=${CHANNEL_TASKS_HOME}"
Environment="CHANNEL_TASKS_INI_PATH=${CHANNEL_TASKS_INI_PATH}"
Environment="CHANNEL_TASKS_PORT=${CHANNEL_TASKS_PORT}"
Environment="CHANNEL_TASKS_STATIC_ROOT=${CHANNEL_TASKS_STATIC_ROOT}"
Environment="CHANNEL_TASKS_STATIC_URI=${CHANNEL_TASKS_STATIC_URI}"
Environment="TASK_ADMIN_USER=${TASK_ADMIN_USER}"
Environment="TASK_ADMIN_EMAIL=${TASK_ADMIN_EMAIL}"

[Install]
WantedBy=multi-user.target
